# scripts/conf/celerybeat.conf
[program:celerybeat]
command=/srv/oscar/project/scripts/celery/beat/start_non_detach.sh
user=oscar
autostart=true
autorestart=true
redirect_stderr=true
