# scripts/conf/gunicorn.conf
[program:gunicorn]
command=/srv/oscar/project/scripts/gunicorn/start_non_detach.sh
user=oscar
autostart=true
autorestart=true
redirect_stderr=true
